# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
# stub